# metal_iperf_instance/roles/iperf/templates/iperf_systemd.service.j2
# /etc/systemd/system/iperf.service
[Unit]
Description=iperf server
After=syslog.target network.target auditd.service

[Service]
ExecStart=/usr/bin/iperf3 -s --logfile=/tmp/iperf_logfile

[Install]
WantedBy=multi-user.target
